# supervisor_conf/job.ini
[program:job]
command=sh -c "echo example job ${MYVAR} ${ANOTHER} && exec sleep 10"
environment=MYVAR="${MYVAR}",ANOTHER="${ANOTHER}"
autostart=false
autorestart=true
;startsecs=1
;startretries=3

